# .crontab
###--- BEGIN DREAMHOST BLOCK
###--- Changes made to this part of the file WILL be destroyed!
# whereami?
MAILTO=""
@reboot sh -c whoami
###--- You can make changes below the next line and they will be preserved!
###--- END DREAMHOST BLOCK

## aggro cron

PHP_PATH='/usr/local/php81/bin/php'
INDEX_PATH='/home/bmxfeed/aggro/current/public/index.php'

## aggro news
## Featured/stream build every 6 minutes
*/6 * * * * $PHP_PATH $INDEX_PATH aggro news

## aggro news/clean
## Clean news featured table every 7 days
0 0 * * 0 $PHP_PATH $INDEX_PATH aggro news-clean

## aggro news/cc
## Clear feed cache files every month (on day 15)
0 0 15 * * $PHP_PATH $INDEX_PATH aggro news-cache

## aggro sweep
## Archive old videos every day (at 2am)
0 2 * * * $PHP_PATH $INDEX_PATH aggro sweep

## aggro vimeo
## Look for new Vimeo videos every 7 minutes
*/7 * * * * $PHP_PATH $INDEX_PATH aggro vimeo

## aggro youtube
## Look for new YouTube videos every 5 minutes
*/5 * * * * $PHP_PATH $INDEX_PATH aggro youtube
